# media/uploads/1522773940clear_backup.cron
#Clear old backups
* * * * * root /backup/clear.php 